# config/systemd/user/omarchy-battery-monitor.service
[Unit]
Description=Omarchy Battery Monitor Check
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/share/omadora/bin/omarchy-battery-monitor
Environment=DISPLAY=:0
